# oidc-mock-renewer/crontab
# min	hour	day	month	weekday	command
*	*	*	*	*	step ca renew --force $CLIENT_CERTS_DIR/site.crt $CLIENT_CERTS_DIR/site.key
